# debian/edgesec.target
[Unit]
Description=EDGESec Target
After=network.target
Conflicts=shutdown.target

[Install]
WantedBy=multi-user.target
